# docker-maintainer/supervisor/judgedaemon1.conf
[program:judgedaemon1]
command=/domjudge/bin/judgedaemon -n 1
numprocs=1
user=domjudge
autostart=true
autorestart=true
