# docker-maintainer/supervisor/judgedaemon0.conf
[program:judgedaemon0]
command=/domjudge/bin/judgedaemon -n 0
numprocs=1
user=domjudge
autostart=true
autorestart=true
