# docker/domserver/supervisor/php.conf
[program:php]
command=php-fpm7.3 -F
numprocs=1
autostart=true
autorestart=true
