# debian/create-cgroups.service
[Unit]
Description=Make sure cgroups exist for domjudge judgedaemon

[Service]
Type=oneshot
ExecStart=/usr/sbin/create_cgroups
RemainAfterExit=true
