# docker-contributor/supervisor/apache2.conf
[program:apache2]
command=pidproxy /var/run/apache2/apache2.pid /bin/bash -c "source /etc/apache2/envvars && apache2ctl -D FOREGROUND"
numprocs=1
autostart=false
autorestart=true
