# docker-contributor/supervisor/judgedaemonextra.conf
[program:judgedaemonextra]
process_name=judgedaemon%(process_num)s
numprocs=2
numprocs_start=2
command=/domjudge/bin/judgedaemon -n %(process_num)s
user=domjudge
autostart=false
autorestart=true
