# docker-contributor/supervisor/judgedaemon.conf
[program:judgedaemon]
process_name=judgedaemon%(process_num)s
numprocs=2
command=PROJECT_DIR/bin/judgedaemon -n %(process_num)s
user=domjudge
autostart=true
autorestart=true
