# docker-contributor/supervisor/nginx.conf
[program:nginx]
command=nginx -g "daemon off;"
numprocs=1
autostart=false
autorestart=true
