# docker-contributor/supervisor/php.conf
[program:php]
command=php-fpm7.4 -F
numprocs=1
autostart=true
autorestart=true
